# scripts/etc/systemd/system/vir.socket
Unit]
Description=wiki socket

[Socket]
ListenStream=/home/user/vir/vir.sock
User=user

[Install]
WantedBy=sockets.target
